# nvidia.service
[Unit]
Description=Set nvidia gpu power limit and clocks on boot
[Service]
ExecStart=/usr/local/sbin/nvidia.sh
[Install]
WantedBy=multi-user.target
